# LXroot/etc/systemd/system/uupd.timer
[Unit]
Description=System Updates enforcer

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
