# Linux/FSRoot/etc/systemd/system/uupd.timer
[Unit]
Description=System Updates enforcer
Wants=network-online.target
After=network-online.target

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
